# update-tools/service/distrobox-update.service
[Unit]
Description=update distrobox

[Service]
Type=oneshot
ExecStart = /home/--your-username-here--/.local/bin/update-distrobox
